# root/etc/s6-overlay/s6-rc.d/proton-natpmp/run
#!/bin/bash
while true
  do
    /app/proton_wg_natpmp.sh
    sleep 45
  done
